# mqtt/rootfs/etc/services.d/mosquitto/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: MQTT Server & Web client
# Runs the MQTT Server
# ==============================================================================
declare -a options

if bashio::config.false 'broker'; then
  bashio::log.debug 'MQTT Server is not enabled...'
  exec sleep 86400
fi

options+=(-c /opt/mosquitto.conf)

if bashio::debug; then
  options+=(-v)
fi

bashio::log.info "Starting mosquitto MQTT Server..."
exec mosquitto "${options[@]}"
