# mqtt/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: MQTT Server & Web client
# Runs the NGINX daemon
# ==============================================================================
if bashio::config.false 'web.enabled'; then
    while true; do
      exec sleep 86400
    done
fi

# Disable HA Authentication if front door is open
if bashio::config.true 'leave_front_door_open'; then
    export DISABLE_HA_AUTHENTICATION=true
fi

bashio::log.info "Starting NGINX for the MQTT client..."
if bashio::config.true 'web.ssl'; then
    options='/etc/nginx/nginx-ssl.conf'
else
    options='/etc/nginx/nginx.conf'
 fi

 exec nginx -c "$options"